# contrib/k70-volume-filter.service
[Unit]
Description=Corsair K70 TKL Champion Volume Wheel Filter
After=multi-user.target

[Service]
Type=simple
ExecStart=k70-volume-filter
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
